# wikihub-insight-master/webapp/simple.conf
[program:wikihub]
command=sudo gunicorn wikihub_worker:app -w 4 -b 0.0.0.0:80

[supervisord]
logfile=/home/ubuntu/supervisord.log
loglevel=debug
user=ubuntu